# scripts/schedule_ceo_briefing.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.digitalfte.ceo-briefing</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/hparacha/DigitalFTE/scripts/weekly_audit.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/hparacha/DigitalFTE</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/hparacha/DigitalFTE/vault/Logs/ceo_briefing.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/hparacha/DigitalFTE/vault/Logs/ceo_briefing_error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>VAULT_PATH</key>
        <string>/Users/hparacha/DigitalFTE/vault</string>
    </dict>
</dict>
</plist>
